# ubuntu-cfg/etc/systemd/system/docker.service.d/overlay.conf
[Service]
# workaround to include default options
EnvironmentFile=/etc/default/docker   
ExecStart=
ExecStart=/usr/bin/dockerd $DOCKER_OPTS
#ExecStart=/usr/bin/docker daemon -H fd:// $DOCKER_OPTS
